# rc.local
#!/bin/sh -e

(sleep 2; /opt/augment00/start.sh)&

exit 0